# docker/root/etc/service/nginx/run
#!/usr/bin/env sh
exec /usr/sbin/nginx -c /www/nginx.conf
